# files/crontab
* * * * * root /usr/local/bin/dnsmasq.sh >> /var/log/cron.log 2>&1

